# rootfs/etc/s6-overlay/s6-rc.d/bazarr/run
#!/usr/bin/with-contenv sh

# change chdir
cd /config

# run with limited permissions
exec s6-setuidgid $PUID:$PGID python3 /app/bazarr.py --debug --no-update --config=/config
